# templates/centos/nad_init.erb
#!/bin/bash
# chkconfig: 2345 95 20
# description: Protect yer NADZ!

PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

case "$1" in
  start)
    NODE_PATH="<%= node['nad']['prefix'] %>/etc/node_modules" node <%= node['nad']['prefix'] %>/sbin/nad -c <%= node['nad']['prefix'] %>/etc/node-agent.d -p <%= @server_address %> &
    echo 'NAD started'
  ;;

  stop)
    NAD_PIDS="$(ps wwwaux | grep node-agent | grep -v grep | awk '{print $2}' | tr "\n" " ")"
    if [ -n "$NAD_PIDS" ] ; then
      kill -9 $NAD_PIDS
    else
      echo 'NAD is not running'
    fi
  ;;

  restart)
    $0 stop
    $0 start
  ;;

  status)
    NAD_PIDS="$(ps wwwaux | grep node-agent | grep -v grep | awk '{print $2}' | tr "\n" " ")"
    if [ -n "$NAD_PIDS" ]; then
      echo "NAD is running on PID(s) $NAD_PIDS!"
      exit 0
    else
      echo 'NAD is not running'
      exit 1
    fi
  ;;

  *)
    echo "Usage: /etc/init.d/nad {start|stop|restart|status}"
    exit 1
  ;;
esac

exit 0
